# tl-service.service
[Unit]
Description=REST-Service which controlles taffic lights by GPIO-pins.

[Service]
ExecStart=tl-service

[Install]
WantedBy=multi-user.target